# fullstack/deploy/systemd/asd-celery-beat.service
[Unit]
Description=ASDSystem Celery Beat
After=network-online.target asd-backend.service
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/srv/asd-system/fullstack/backend
EnvironmentFile=/srv/asd-system/fullstack/deploy/env/backend.production.env
ExecStart=/srv/asd-system/fullstack/backend/.venv/bin/celery -A config.celery:app beat --loglevel=info --pidfile=
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=5
KillMode=mixed
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
